# server/prism-server.service
[Unit]
Description=PRISM Server (Infrastructure Management Hub)
After=network.target

[Service]
Type=simple
User=prism
WorkingDirectory=/opt/prism/server
ExecStart=/opt/prism/server/prism-server
Restart=always
RestartSec=5
Environment=PORT=65432

[Install]
WantedBy=multi-user.target
